# .config/systemd/user/lxpolkit.service
[Unit]
Description=PolicyKit agent for LXDE

[Service]
Type=simple
ExecStart=/usr/bin/lxpolkit

[Install]
WantedBy=default.target
