# .config/systemd/user/lxpolkit.service
[Unit]
Description=PolicyKit agent for LXDE

[Service]
Type=forking
ExecStart=/usr/bin/lxpolkit

[Install]
WantedBy=multi-user.target