# .config/autostart/fcitx5.desktop
[Desktop Entry]
Name=fcitx5
Type=Application
Exec=/usr/bin/fcitx5
