# .config/autostart/otd-daemon.desktop
[Desktop Entry]
Name=otd-daemon
Comment=User mode FOSS tablet driver
Type=Application
Exec=/usr/bin/otd-daemon
